# .config/systemd/user/wallpaper.service
[Unit]
Description=Change wallpaper
Requires=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wallpaper-shuffle

